# ansible/roles/bookstack/files/bookstack.service
[Unit]
Description=BookStack - A platform to create documentation/wiki content built with PHP & Laravel
Documentation=https://www.bookstackapp.com/

Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/etc/containers-compose/bookstack
ExecStart=/usr/bin/docker compose up -d --remove-orphans --build --force-recreate
ExecStop=/usr/bin/docker compose down

[Install]
WantedBy=multi-user.target
